# telco-ran/configuration/extra-manifests-builder/01-container-mount-ns-and-kubelet-conf/container-mount-namespace.service
[Unit]
Description=Manages a mount namespace that both kubelet and crio can use to share their container-specific mounts

[Service]
Type=oneshot
RemainAfterExit=yes
RuntimeDirectory=container-mount-namespace
Environment=RUNTIME_DIRECTORY=%t/container-mount-namespace
Environment=BIND_POINT=%t/container-mount-namespace/mnt
ExecStartPre=bash -c "findmnt ${RUNTIME_DIRECTORY} || mount --make-unbindable --bind ${RUNTIME_DIRECTORY} ${RUNTIME_DIRECTORY}"
ExecStartPre=touch ${BIND_POINT}
ExecStart=unshare --mount=${BIND_POINT} --propagation slave mount --make-rshared /
ExecStop=umount -R ${RUNTIME_DIRECTORY}

# Initialize /run/netns as a shared mount point to prevent shadowing race condition
# This must happen before CRI-O starts to ensure all namespace bind mounts have
# the correct parent mount. Without this, 'ip netns add' (called later by network
# components) would create the mount point and shadow any earlier namespace mounts
# created by CRI-O's pinns, causing pods to fail with 'setns: Invalid argument'.
# Creating and deleting a temporary netns retains the side-effect of creating the
# /run/netns mountpoint (but base exit code solely on whether the mountpoint is
# present)
# Reference: OCPBUGS-83562
ExecStartPost=bash -c "ip netns add kubensmnt && ip netns delete kubensmnt || mountpoint -q /run/netns"
